# eos-updater-avahi/eos-updater-avahi.service.in
[Unit]
Description=Endless OS Avahi Advertisement Updater
Documentation=man:eos-updater-avahi(8)

[Service]
Type=simple
ExecStart=@libexecdir@/eos-updater-avahi

# Sandboxing
# FIXME: Enable more of these options once we have systemd > 229
CapabilityBoundingSet=
Environment=GIO_USE_VFS=local
Environment=GVFS_DISABLE_FUSE=1
Environment=GVFS_REMOTE_VOLUME_MONITOR_IGNORE=1
#MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateNetwork=yes
PrivateTmp=yes
#PrivateUsers=yes
#ProtectControlGroups=yes
ProtectHome=yes
#ProtectKernelModules=yes
#ProtectKernelTunables=yes
#ProtectSystem=strict
ProtectSystem=true
#ReadWritePaths=/etc/avahi/services/eos-updater.service
#RestrictRealtime=yes
#SystemCallErrorNumber=EPERM
#SystemCallFilter=@basic-io @io-event @ipc

[Install]
Also=eos-update-server.service eos-updater-avahi.path
