# eos-autoupdater/eos-autoupdater-payg.timer.in
[Unit]
Description=Endless OS Automatic Update Timer (PAYG)
Documentation=man:eos-autoupdater(8)
ConditionKernelCommandLine=!endless.live_boot
ConditionKernelCommandLine=ostree
ConditionKernelCommandLine=eospayg

[Timer]
OnBootSec=0m
OnUnitInactiveSec=1h
RandomizedDelaySec=30min
Unit=eos-autoupdater.service

[Install]
WantedBy=multi-user.target
