# eos-autoupdater/eos-autoupdater.timer.in
# Copyright 2018, 2019 Endless Mobile, Inc.
# SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Endless OS Automatic Update Timer
Documentation=man:eos-autoupdater(8)
ConditionKernelCommandLine=!endless.live_boot
ConditionKernelCommandLine=ostree

[Timer]
# After 3 minutes the network connection has likely been set up, but we don't
# have a hard requirement; offline updates need to work too
OnActiveSec=3m
OnUnitInactiveSec=1h
RandomizedDelaySec=30min

[Install]
WantedBy=multi-user.target
