# eos-updater-flatpak-installer/eos-updater-flatpak-installer.service.in
# Copyright 2017, 2018, 2019 Endless Mobile, Inc.
# SPDX-License-Identifier: LGPL-2.1-or-later

[Unit]
Description=Endless OS Post-Boot Flatpak Installer
Documentation=man:eos-updater-flatpak-installer(8)
# /home is a symlink to /var/home; /var/home is a symlink to /sysroot/home. The
# second symlink is created by systemd-tmpfiles. Since we use ProtectHome=yes,
# we must explicitly order this unit after tmpfiles are created.
Requires=local-fs.target systemd-tmpfiles-setup.service
After=local-fs.target systemd-tmpfiles-setup.service
ConditionKernelCommandLine=!eos-updater-disable
DefaultDependencies=no
Conflicts=shutdown.target

# Only on updates
Before=multi-user.target systemd-update-done.service
ConditionNeedsUpdate=/var

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=@libexecdir@/eos-updater-flatpak-installer
Restart=no

# Flatpak checks parental controls at deploy time. In order to do this, it
# needs to talk to accountsservice on the system bus, neither of which are
# running when this job runs.
Environment=FLATPAK_SKIP_PARENTAL_CONTROLS_NO_SYSTEM_BUS=1

# Sandboxing
# flatpak triggers use bwrap which requires net/sys admin and chroot
CapabilityBoundingSet=CAP_NET_ADMIN CAP_SYS_ADMIN CAP_SYS_CHROOT
Environment=GIO_USE_VFS=local
Environment=GVFS_DISABLE_FUSE=1
Environment=GVFS_REMOTE_VOLUME_MONITOR_IGNORE=1
Environment=GSETTINGS_BACKEND=memory
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateNetwork=yes
PrivateTmp=yes
PrivateUsers=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
# bwrap also mounts /proc
ProtectKernelTunables=no
ProtectSystem=no
RestrictAddressFamilies=AF_NETLINK AF_UNIX
RestrictRealtime=yes
SystemCallArchitectures=native
SystemCallErrorNumber=EPERM
# @network-io is required for logging to the journal to work
# @privileged and @chown are required for certain ostree operations
# @mount is required for bwrap
SystemCallFilter=~@clock @cpu-emulation @debug @keyring @module @obsolete @raw-io @resources

[Install]
WantedBy=multi-user.target
