# eos-updater-flatpak-installer/eos-updater-flatpak-installer.service.in
[Unit]
Description=Endless OS Post-Boot Flatpak Installer
Wants=local-fs.target
After=local-fs.target
ConditionKernelCommandLine=!eos-updater-disable
DefaultDependencies=no

# Only on updates
Before=multi-user.target systemd-update-done.service
ConditionNeedsUpdate=/var

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=@libexecdir@/eos-updater-flatpak-installer

# Sandboxing
CapabilityBoundingSet=
Environment=GIO_USE_VFS=local
Environment=GVFS_DISABLE_FUSE=1
Environment=GVFS_REMOTE_VOLUME_MONITOR_IGNORE=1
Environment=GSETTINGS_BACKEND=memory
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateNetwork=yes
PrivateTmp=yes
PrivateUsers=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=no
RestrictAddressFamilies=AF_UNIX
RestrictRealtime=yes
SystemCallErrorNumber=EPERM
SystemCallFilter=@basic-io @io-event @ipc @network-io

[Install]
WantedBy=multi-user.target
