# data/eos-updater.service.in
[Unit]
Description=Endless OS Updater
Documentation=man:eos-updater(8)
Wants=network.target

[Service]
ExecStart=@bindir@/eos-updater
Type=dbus
BusName=com.endlessm.Updater

# Sandboxing
# FIXME: Enable more of these options once we have systemd > 229
#CapabilityBoundingSet=(some list of capabilities)
Environment=GIO_USE_VFS=local
Environment=GVFS_DISABLE_FUSE=1
Environment=GVFS_REMOTE_VOLUME_MONITOR_IGNORE=1
#MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateNetwork=no
PrivateTmp=yes
#PrivateUsers=yes
#ProtectControlGroups=yes
ProtectHome=yes
#ProtectKernelModules=yes
#ProtectKernelTunables=yes
#ProtectSystem=strict
ProtectSystem=no
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
#RestrictRealtime=yes
#SystemCallErrorNumber=EPERM
#SystemCallFilter=@basic-io @io-event @ipc @network-io
