# data/eos-auto-updates.timer.in
[Unit]
Description=EndlessOS Automatic Update Timer
Wants=network.target

[Timer]
OnBootSec=15m
OnUnitInactiveSec=1h

[Install]
WantedBy=multi-user.target
