# data/eos-update-server.socket.in
[Unit]
Description=Socket for Endless OS Automatic Updater

[Socket]
ListenStream=@server_port@

[Install]
WantedBy=sockets.target
