# data/eos-release-upgrade.service.in
[Unit]
Description=EndlessOS Release Upgrader
Requires=local-fs.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/eos-release-upgrade
