# data/eos-autoupdater.timer.in
[Unit]
Description=Endless OS Automatic Update Timer
Documentation=man:eos-autoupdater(8)
Wants=network.target
ConditionKernelCommandLine=!endless.live_boot

[Timer]
OnBootSec=15m
OnUnitInactiveSec=1h
RandomizedDelaySec=30min

[Install]
WantedBy=multi-user.target
