# data/eos-updater.service.in
[Unit]
Description=EndlessOS Updater
Wants=network.target

[Service]
Type=simple
ExecStart=@sbindir@/eos-updater
