# data/eos-update-server.socket.in
[Unit]
Description=Socket for EndlessOS Automatic Updater

[Socket]
ListenStream=@server_port@
