# data/eos-update-server.service.in
[Unit]
Description=Endless OS Local Update Server
Documentation=man:eos-update-server(8)
Wants=network.target

# Only run when Avahi is advertising the port.
ConditionPathExists=@sysconfdir@/avahi/services/eos-updater.service

[Service]
Type=simple
ExecStart=@libexecdir@/eos-update-server

# Deprioritise the process to reduce UI impact if clients are downloading from
# this machine. Numbers are chosen fairly arbitrarily.
Nice=15
IOSchedulingClass=idle

# Sandboxing
# FIXME: Enable more of these options once we have systemd > 229
CapabilityBoundingSet=CAP_DAC_READ_SEARCH
Environment=GIO_USE_VFS=local
Environment=GVFS_DISABLE_FUSE=1
Environment=GVFS_REMOTE_VOLUME_MONITOR_IGNORE=1
Environment=GSETTINGS_BACKEND=memory
#MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
PrivateDevices=yes
PrivateNetwork=no
PrivateTmp=yes
#PrivateUsers=yes
#ProtectControlGroups=yes
ProtectHome=yes
#ProtectKernelModules=yes
#ProtectKernelTunables=yes
#ProtectSystem=strict
ProtectSystem=full
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
#RestrictRealtime=yes
#SystemCallErrorNumber=EPERM
#SystemCallFilter=@basic-io @io-event @ipc @network-io

[Install]
Also=eos-update-server.socket
