# data/eos-auto-updates.service.in
[Unit]
Description=EndlessOS Updater
Wants=network.target

[Service]
Type=simple
ExecStart=@libexecdir@/eos-auto-updates
