# data/eos-update-server.socket.in
[Unit]
Description=Socket for Endless OS Automatic Updater
Documentation=man:eos-update-server(8)
ConditionPathExists=@sysconfdir@/avahi/services/eos-updater.service

[Socket]
ListenStream=@server_port@

[Install]
WantedBy=sockets.target
