# data/eos-autoupdater.timer.in
[Unit]
Description=Endless OS Automatic Update Timer
Wants=network.target

[Timer]
OnBootSec=15m
OnUnitInactiveSec=1h

[Install]
WantedBy=multi-user.target
Alias=eos-updater.timer
