# data/eos-update-server.service.in
[Unit]
Description=EndlessOS Local Update Server
Wants=network.target

[Service]
Type=simple
ExecStart=@libexecdir@/eos-update-server
