# data/eos-updater.service.in
[Unit]
Description=EndlessOS Updater
ConditionKernelCommandLine=ostree
Wants=network.target

[Service]
Type=simple
ExecStart=@libexecdir@/eos-updater
